# deploy/miao.conf

description "Miao Server"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

script
exec python /root/workspace/Miao/manage.py runfcgi host=0.0.0.0 port=8056 pidfile=/var/run/miao_server.pid daemonize=true

end script

post-start script
end script